# shell-tools/share/pxi-keyword-insight/playwright-api.service
[Unit]
Description=Playwright API (keyword-insight /navigate runtime)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/playwright/docker
Environment=PORT=3000
Environment=HOME=/root
Environment=PROXY_SERVER=
Environment=PROXY_USERNAME=
Environment=PROXY_PASSWORD=
ExecStart=/usr/bin/node /opt/playwright/docker/server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
